# .config/systemd/user/dl_monthly.service


[Unit]
Description=/usr/bin/fish -c dl_monthly
Requires=mnt-d.mount

[Service]
RemainAfterExit=yes
ExecStart=
ExecStart="/usr/bin/fish" "-c" "dl_monthly"
