# .config/systemd/user/tabs-cine.timer


[Unit]
Description=/usr/bin/fish -c cine

[Timer]
Persistent=yes
OnCalendar=*-*-* 9:55
RemainAfterElapse=yes
